# crontab
0	*	*	*	*	~/.config/bin/user-env.cmd ~/.config/bin/update_crontab.py
