# crontab
0	*	*	*	*	~/.config/bin/crontab-env.cmd ~/.config/bin/update_crontab.py
